# launchd/templates/com.dailytech.site.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dailytech.site.tunnel</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-lc</string>
      <string>cd "__PROJECT_DIR__" &amp;&amp; exec zsh scripts/run-cloudflare-tunnel.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/cloudflare-tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/cloudflare-tunnel.err.log</string>
  </dict>
</plist>
