# launchd/templates/com.dailytech.site.refresh.afternoon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dailytech.site.refresh.afternoon</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-lc</string>
      <string>cd "__PROJECT_DIR__" &amp;&amp; exec zsh scripts/run-refresh-check.sh afternoon</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>20</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/refresh-afternoon.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/refresh-afternoon.err.log</string>
  </dict>
</plist>
